# mvsg.cron
*/1 * * * * root sh /usr/local/mvsg/mvsg.sh > /dev/null 2>&1